# services/execution-bridge/olymp-bridge.service
[Unit]
Description=OLYMP Execution Bridge
Documentation=https://github.com/aicallyu/olympus
After=network.target

[Service]
Type=simple
User=onioko
WorkingDirectory=/home/onioko/olympus/services/execution-bridge
ExecStart=/home/onioko/.nvm/versions/node/v24.13.0/bin/node /home/onioko/olympus/services/execution-bridge/dist/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=olymp-bridge

# Environment
EnvironmentFile=/home/onioko/olympus/services/execution-bridge/.env

# Safety
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
